# docker/etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv sh
printenv > /etc/environment
cron -f